# data/etc/s6-overlay/s6-rc.d/provideHooks/up
foreground { echo "(custom-hook-provider): Copying hooks to volume /hooks for use by certbot..." }
foreground { cp /data/hooks/auth-hook.sh /hooks/ }