# data/etc/s6-overlay/s6-rc.d/provideHooks/up
foreground { echo "copying hooks to volume /hooks for use by certbot..." }
cp /data/hooks/* /hooks/