# data/etc/s6-overlay/s6-rc.d/makeExecutable/up
foreground { echo "making everything in /make-executable executable..." }
chmod +x /make-executable/* 2>/dev/null
